# .config/systemd/user/xdg-desktop-portal.service.d/restart.conf
[Service]
# Auto-restart on crash instead of staying dead
Restart=on-failure
RestartSec=3
