# .config/systemd/user/dbus-broker.service.d/limits.conf
[Service]
# Increase dbus limits to prevent peer disconnections under memory pressure
# These are set via busconfig, not CLI flags
LimitNOFILE=8192
LimitMSGQUEUE=8388608
